# deploy/systemd/argos-agent.service
[Unit]
Description=ARGOS Agent (recolector de telemetría hacia el servidor ARGOS)
After=network-online.target argos-server.service
Wants=network-online.target

[Service]
Type=simple
User=argos
Group=argos
WorkingDirectory=/opt/argos
ExecStart=/opt/argos/.venv/bin/python -m argos.agent
Environment=ARGOS_DATA_DIR=/var/lib/argos
Environment=ARGOS_SERVER_HOST=127.0.0.1
Environment=ARGOS_SERVER_PORT=8000
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
